# docker/laravel-echo/supervisor-codekoalas-conf.ini
[program:laravel-echo-server]
command = /opt/var/init.sh
startsecs = 0
autorestart = true
startretries = 3
